# project/ansible/back_role/templates/apache_exporter.service
[Unit]
Description=Apache exporter Service
After=network.target  

[Service]
ExecStart=/usr/bin/apache_exporter

[Install]
WantedBy=multi-user.target  