# les27/templates/borg-backup.timer
[Unit]
Description=Borg Backup

[Timer]
OnBootSec=10s
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
