# project/ansible/node_exp_role/templates/node_exporter.service
[Unit]
Description=Node exporter Service
After=network.target
[Service]
ExecStart=/usr/bin/node_exporter
[Install]
WantedBy=multi-user.target