# project/ansible/db_backup_role/templates/backup.service
[Unit]
Description=Mysql Backup

[Service]
Type=oneshot

# Создание бэкапа
ExecStart=/bin/bash /backup/run_backup.sh

