# project/ansible/mon_role/templates/alertmanager.service
[Unit]
Description=Alertmanager Service
After=network.target
[Service]
ExecStart=/usr/local/bin/alertmanager
[Install]
WantedBy=multi-user.target